# supervisor_config.ini
[supervisord]
logfile=./logs/supervisord.log ; (default $CWD/supervisord.log)
logfile_maxbytes=50MB   ; (default 50MB)
logfile_backups=10      ; (default 10)
loglevel=info           ; (default info)
pidfile=supervisord_asr.pid  ; (default supervisord.pid)

[unix_http_server]
file=/var/run/supervisor_asr.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)

[supervisorctl]
serverurl=unix:///var/run/supervisor_asr.sock ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:asr_api]
command=uvicorn app.main:app --host=0.0.0.0 --port=8000 --workers=4
directory=./
autostart=true
autorestart=true
stdout_logfile=./logs/api.log
stderr_logfile=./logs/api_error.log
environment=PYTHONPATH=/usr/local/Caskroom/miniforge/base/envs/videotrans,PYTHONUNBUFFERED=1,PYTHONASYNCIODEBUG=1

[program:asr_celery]
command=python3 celery_worker.py
directory=./
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./logs/celery.log
stderr_logfile=./logs/celery_error.log
environment=PYTHONPATH=/usr/local/Caskroom/miniforge/base/envs/videotrans,PYTHONUNBUFFERED=1,PYTHONASYNCIODEBUG=1

[program:asr_flower]
command=python3 celery_flower.py
directory=./
autostart=true
autorestart=true
stdout_logfile=./logs/flower.log
stderr_logfile=./logs/flower_error.log
environment=PYTHONPATH=/usr/local/Caskroom/miniforge/base/envs/videotrans,PYTHONUNBUFFERED=1,PYTHONASYNCIODEBUG=1

[group:asr_service]
programs=asr_api,asr_celery,asr_flower
priority=999 