# supervisor_config.ini
[supervisord]
logfile=./logs/supervisord.log ; (default $CWD/supervisord.log)
logfile_maxbytes=50MB   ; (default 50MB)
logfile_backups=10      ; (default 10)
loglevel=info           ; (default info)
pidfile=supervisord_asr.pid  ; (default supervisord.pid)

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9001        ; ip_address:port specifier, *:port for all iface
username=user              ; default is no username (open server)
password=123               ; default is no password (open server)

[unix_http_server]
file=/var/run/supervisor_asr.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)

[supervisorctl]
serverurl=unix:///var/run/supervisor_asr.sock ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:asr_api]
command=/home/ubuntu/miniforge3/envs/asr/bin/uvicorn app.main:app --host=0.0.0.0 --port=8000 --workers=4
directory=./
autostart=true
autorestart=true
stdout_logfile=./logs/api.log
stderr_logfile=./logs/api_error.log
environment=PYTHONUNBUFFERED=1,PYTHONASYNCIODEBUG=1

[program:asr_celery]
command=/home/ubuntu/miniforge3/envs/asr/bin/python celery_worker.py
directory=./
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./logs/celery.log
stderr_logfile=./logs/celery_error.log
environment=PYTHONUNBUFFERED=1,PYTHONASYNCIODEBUG=1

[program:asr_celery_beat]
command=/home/ubuntu/miniforge3/envs/asr/bin/celery -A app.core.celery:celery_app beat --loglevel=info
directory=./
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=./logs/celery_beat.log
stderr_logfile=./logs/celery_beat_error.log
environment=PYTHONUNBUFFERED=1,PYTHONASYNCIODEBUG=1

[group:asr_service]
programs=asr_api,asr_celery,asr_celery_beat
priority=999 